# server/systemd/zerotier-watchdog.service
[Unit]
Description=ZeroTier Watchdog - reinicia se OFFLINE
After=zerotier-one.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zerotier-watchdog.sh
