# server/profinet/systemd/pnet-cmd.service
[Unit]
Description=Executar comandos PROFINET do PLC

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pnet-cmd-handler.sh
