# server/profinet/systemd/pnet-status.service
[Unit]
Description=Actualiza estado do sistema para PROFINET p-net

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pnet-status-update.sh
